# benchmarks/cluster-setup/mongochain.service
# Copy to /usr/lib/systemd/system/mongochain.service
[Unit]
Description=Mongochain

[Service]
User=root
WorkingDirectory=/home/akka
ExecStart=/home/akka/workspace/mongochain/benchmarks/cluster-setup/update-start.sh

[Install]
WantedBy=multi-user.target
